# scripts/autostart/startup_app.desktop
[Desktop Entry]
Type=Application
Name=custom_startup
Exec=/home/rpiarcade/startup.sh
Comment=Runs backend on startup